# templates/systemd/docker-registry-garbage-collect.service.j2
[Unit]
Description=Performs container registry garbage collection

[Service]
Type=oneshot
ExecStart={{ docker_registry_bin_path }}/garbage-collect
